# target/linux/hisilicon/base-files/etc/init.d/done
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2006 OpenWrt.org
#

START=95
boot() {
	mount_root done
	rm -f /sysupgrade.tgz

	# process user commands
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}

	# gpio init
	. /etc/diag.sh
	gpio_init

	# set leds to normal state
	. /etc/diag.sh
	set_state done

	# run cmd from u-boot - idea by Maxim - realisation by Igor ;)
	linux_cmd=$(fw_printenv -n linux_cmd) ; [ -n "${linux_cmd}" ] && eval "${linux_cmd}" >/dev/null 2>&1
}
