# target/linux/hi35xx/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/hi35xx.sh

board=$(hi35xx_board_name)

case "$board" in
hi3516cv100)
	echo "Device hi3516cv100 found - 01_leds" | logger
	;;
hi3516cv200)
	echo "Device hi3516cv200 found - 01_leds" | logger
	;;
hi3518ev100)
	echo "Device hi3518ev100 found - 01_leds" | logger
	;;
hi3518ev200)
	echo "Device hi3518ev200 found - 01_leds" | logger
	;;

jvt_s130h18v)
	echo "Device jvt_s130h18v found - 01_leds" | logger
	;;
jvt_s135h18v)
	echo "Device jvt_s135h18v found - 01_leds" | logger
	;;
jvt_s323h16v)
	echo "Device jvt_s323h16v found - 01_leds" | logger
	;;
esac

ucidef_commit_leds

exit 0
