# target/linux/hi35xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2014 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/hi35xx.sh

ucidef_set_interface_loopback

board=$(hi35xx_board_name)

case "$board" in
*)
	local mac
	mac=$(fw_printenv -n ethaddr)
	[ -n "$mac" ] && ucidef_set_interface_macaddr "lan" "$mac"
	echo "Device $board found - 02_network" | logger
	;;
esac

uci commit network

exit 0
