# target/linux/hisilicon/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/hisilicon.sh

board=$(hisilicon_board_name)

case "$board" in
hi3518ev100)
	#ucidef_set_led_usbdev "usb" "USB" "dlan-usb-extender:green:usb" "2-1"
	echo "Device hi3518ev100 found - 01_leds" | logger
	;;
hi3518ev100)
	#ucidef_set_led_netdev "eth" "eth" "rbt-832:red:ethernet" "eth0"
	echo "Device hi3518ev200 found - 01_leds" | logger
	;;
esac

ucidef_commit_leds

exit 0
