# target/linux/hisilicon/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/hisilicon.sh

board=$(hisilicon_board_name)

case "$board" in
hi3516cv100)
	echo "Device ${board} found - 01_leds" | logger
	;;
hi3516cv200)
	echo "Device ${board} found - 01_leds" | logger
	;;
hi3518ev100)
	echo "Device ${board} found - 01_leds" | logger
	;;
hi3518ev200)
	echo "Device ${board} found - 01_leds" | logger
	;;

jvt_s130h18v)
	echo "Device ${board} found - 01_leds" | logger
	;;
jvt_s135h18v)
	echo "Device ${board} found - 01_leds" | logger
	;;
jvt_s323h16v)
	echo "Device ${board} found - 01_leds" | logger
	;;
esac

ucidef_commit_leds

exit 0
